# mirrorlist_clientOMA.ini
[program:mirrorlist_server]
command=python  /var/lib/mirrormanager/pyvirt/mirrormanager/usr/share/mirrormanager/mirrorlist_server.py
priority=990
username=mirrormanager
environment=PATH=/var/lib/mirrormanager/pyvirt/mirrormanager/bin/
stdout_logfile=/tmp/%(program_name)s.stdout
stderr_logfile=/tmp/%(program_name)s.stderr
autostart=true
autorestart=true
