# mirrormanagerOMA.ini
[program:mirrormanager]
command=/var/lib/mirrormanager/pyvirt/mirrormanager/bin/uwsgi  --emperor /etc/uwsgi/vassals/
priority=700
username=mirrormanager
stdout_logfile=/tmp/%(program_name)s.stdout
stderr_logfile=/tmp/%(program_name)s.stderr
autostart=true
autorestart=true
